# scripts/engram-sync-daemon.service
[Unit]
Description=Engram Memory Sync Daemon (Framework-SDD)
After=network.target

[Service]
Type=simple
User=%u
# Secretos: ~/.config/framework-sdd/engram-daemon.env (ver config/engram-daemon.env.example)
EnvironmentFile=-%h/.config/framework-sdd/engram-daemon.env
# Ajustar la ruta si el clon no está en esta ruta (o usar symlink):
ExecStart=%h/Documents/Good4D/Framework-SDD/scripts/engram-sync-daemon.sh start
ExecStop=%h/Documents/Good4D/Framework-SDD/scripts/engram-sync-daemon.sh stop
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
